# etc/systemd/system/zfs-auto-snapshot-daily.service
[Unit]
Description=ZFS Auto Snapshot Daily

[Service]
Type=oneshot
ExecStart=/usr/sbin/zfs-auto-snapshot --quiet --syslog --label=daily --keep=31
